# docker/etc/systemd/system/vitastor.target
[Unit]
Description=vitastor target
[Install]
WantedBy=multi-user.target
